# deploy/ngrok.service
[Unit]
Description=NGROK Tunnel for Anomaly Detector
After=network.target

[Service]
Type=simple
User=pi
ExecStart=/usr/local/bin/ngrok http --domain=your-domain.ngrok.io 8000
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
